# deploy/systemd/mithra-pdf-parser.service
[Unit]
Description=Mithra isolated PDF parser
Requires=mithra-pdf-parser.socket
After=mithra-pdf-parser.socket

[Service]
Type=simple
ExecStart=/usr/local/bin/mithra pdf-parser
User=mithra-pdf
Group=mithra
UMask=0077
NoNewPrivileges=true
PrivateNetwork=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
InaccessiblePaths=/var/lib/mithra/mithra.sqlite3 /var/lib/mithra/sources /etc/mithra/credentials
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LimitNOFILE=64
MemoryMax=192M
CPUQuota=50%
TimeoutStartSec=15s
TimeoutStopSec=5s
RuntimeMaxSec=5min
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
