# deploy/systemd/mithra-backup.timer
[Unit]
Description=Daily Mithra backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=20m

[Install]
WantedBy=timers.target
